# coopr-provisioner/worker/plugins/automators/chef_solo_automator/resources/cookbooks/modules/files/default/modules-load.conf
# modules-load - set /etc/modules and load modules from /etc/modules-load.d/*.conf

description	"set /etc/modules and load modules from /etc/modules-load.d/*.conf"

start on virtual-filesystems or static-network-up

task
script
	if [ -r /etc/modules-load.d/header ] ; then cat /etc/modules-load.d/header > /etc/modules ; else  > /etc/modules ; fi
	if ls /etc/modules-load.d/*.conf > /dev/null 2>&1 ; then cat /etc/modules-load.d/*.conf >> /etc/modules ; fi
end script
